# config/files/my-script/files/my_script
#!/bin/sh /etc/rc.common

START=99

start() {
    sleep 10
    tc qdisc replace dev pppoe-wan root fq_pie
    tc qdisc replace dev eth1 root fq_pie
    tc qdisc replace dev eth0 root fq_pie
}
